# supervisor/supervisor-apache.conf
[program:apache-server]
command=/usr/sbin/apache2ctl -D FOREGROUND
